# upstart/jwt.conf
#!upstart
description "node.js jwt upstart"
author      "sperreault@alesium.net"

start on startup
stop on shutdown

script
    export HOME="/vol0/jwt"

    echo $$ > /var/run/jwt.pid
    exec sudo -u jwt /usr/bin/node $HOME/app.js >> /var/log/jwt.log 2>&1
end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/jwt.log
end script

pre-stop script
    rm /var/run/jwt.pid
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/jwt.log
end script
